# crontab
0 */6 * * * python /pnrStatus.py >> /tmp/logfile.log
